# etc/thingler.cron
#
# crontab for couchdb user
#
#   install with `sudo crontab thingler.cron -u couchdb`
#

# Hourly & Daily backups of the thingler database
00 * * * * cp /var/lib/couchdb/thingler.couch /usr/local/var/couchdb/thingler-backup.hourly.couch
00 0 * * * cp /var/lib/couchdb/thingler.couch /usr/local/var/couchdb/thingler-backup.daily.couch

# Compact database every Sunday at 6AM
00 6 * * 0 curl -H "Content-Type: application/json" -X POST http://0.0.0.0:5984/thingler/_compact
